# bluehood.service
[Unit]
Description=Bluehood Bluetooth Monitor
After=bluetooth.target network.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m bluehood.daemon
Restart=on-failure
RestartSec=10
Environment=BLUEHOOD_DATA_DIR=/var/lib/bluehood

[Install]
WantedBy=multi-user.target
